# conf.d/websockify.conf
[program:websockify]
command=/system/conf.d/websockify.sh
autorestart=true
